# systemd/dispatch-service.service
# dispatch-service - internal processor + callback to ride-api (port 3003).
# Internal-only: binds to 127.0.0.1 so it is unreachable from off-box.
#
# Install:  sudo cp systemd/dispatch-service.service /etc/systemd/system/
#           sudo systemctl daemon-reload
#           sudo systemctl enable --now dispatch-service
[Unit]
Description=dispatch-service (internal processor + callback, port 3003)
Documentation=https://github.com/nebyathhailu/production-service-env
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=serviceenv
Group=serviceenv
WorkingDirectory=/opt/service-env/services/dispatch-service
ExecStart=/opt/service-env/venv/bin/python /opt/service-env/services/dispatch-service/app.py

# Configuration (matches in-code defaults; overridable here).
Environment=DISPATCH_SERVICE_PORT=3003
Environment=BIND_HOST=127.0.0.1
Environment=RIDE_API_URL=http://ride-api.internal:3001
Environment=DOWNSTREAM_TIMEOUT=5

Restart=on-failure
RestartSec=2

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dispatch-service

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
