# systemd/service-a.service
# Service A - public-facing entry point behind Nginx (port 3001).
# Binds to 127.0.0.1; only Nginx (port 80) is public.
#
# DEPENDENCY MANAGEMENT (assignment requirement):
#   Service A must not start before Service B and Service C, and must not
#   become operational until both are actually reachable.
#     - After=        -> ordering: A starts after B and C units have started.
#     - Wants=        -> pulls B/C in as a startup preference (paired with
#                        After= above). Deliberately not Requires=: Requires=
#                        propagates *deactivation* too, so if B or C is later
#                        stopped/crashes, A would be torn down with it - taking
#                        the public entry point down along with one internal
#                        dependency. Wants= keeps the startup ordering without
#                        that runtime coupling, so A stays up and its own
#                        error handling deals with a downstream outage.
#     - ExecStartPre= -> readiness gate: A blocks until B and C answer their
#                        /health endpoints (ordering alone only guarantees the
#                        processes were launched, not that they are listening).
#
# Install:  sudo cp systemd/service-a.service /etc/systemd/system/
#           sudo systemctl daemon-reload
#           sudo systemctl enable --now service-a
[Unit]
Description=Service A (public entry point, port 3001)
Documentation=https://github.com/nebyathhailu/production-service-env
After=network-online.target service-b.service service-c.service
Wants=network-online.target
Wants=service-b.service service-c.service
# Allow several quick retries (e.g. while deps are still coming up) before
# systemd gives up, instead of tripping the default start-limit throttle.
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
User=serviceenv
Group=serviceenv
WorkingDirectory=/opt/service-env/services/service-a

# Readiness gate: wait (up to ~30s) for B and C to actually serve /health
# before launching A. If they never come up, A fails and is retried.
ExecStartPre=/usr/bin/env bash -c 'for i in $(seq 1 30); do \
    curl -sf http://service-b.internal:3002/health >/dev/null \
    && curl -sf http://service-c.internal:3003/health >/dev/null \
    && exit 0; \
    sleep 1; \
  done; echo "dependencies (service-b/service-c) not ready" >&2; exit 1'
ExecStart=/opt/service-env/venv/bin/python /opt/service-env/services/service-a/app.py

# Configuration (matches in-code defaults; overridable here).
Environment=SERVICE_A_PORT=3001
Environment=BIND_HOST=127.0.0.1
Environment=SERVICE_B_URL=http://service-b.internal:3002
Environment=DOWNSTREAM_TIMEOUT=5

Restart=on-failure
RestartSec=2

StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-a

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
