# systemd/service-c.service
# Service C - internal processor + callback to Service A (port 3003).
# Internal-only: binds to 127.0.0.1 so it is unreachable from off-box.
#
# Install:  sudo cp systemd/service-c.service /etc/systemd/system/
#           sudo systemctl daemon-reload
#           sudo systemctl enable --now service-c
[Unit]
Description=Service C (internal processor + callback, port 3003)
Documentation=https://github.com/nebyathhailu/production-service-env
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=serviceenv
Group=serviceenv
WorkingDirectory=/opt/service-env/services/service-c
ExecStart=/opt/service-env/venv/bin/python /opt/service-env/services/service-c/app.py

# Configuration (matches in-code defaults; overridable here).
Environment=SERVICE_C_PORT=3003
Environment=BIND_HOST=127.0.0.1
Environment=SERVICE_A_URL=http://service-a.internal:3001
Environment=DOWNSTREAM_TIMEOUT=5

Restart=on-failure
RestartSec=2

StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-c

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
